# files/etc/rc.local
#!/bin/sh

# 修复 R2C Plus (eth0) 网卡
/usr/sbin/ethtool -K eth0 generic-receive-offload off
/usr/sbin/ethtool -K eth0 tx-scatter-gather off

exit 0
